# systemd/bt-hub.service
[Unit]
Description=Pi BT Hub - Bluetooth Management Web UI
Documentation=https://github.com/hemna/pi-bt-hub
After=bluetooth.target dbus.service network.target
Requires=bluetooth.target dbus.service

[Service]
Type=simple

# Update paths to match your installation (replace /home/pi with your home directory)
WorkingDirectory=/home/pi/pi-bt-hub
ExecStart=/home/pi/pi-bt-hub/.venv/bin/uvicorn bt_hub.main:app --host 0.0.0.0 --port 8080

Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Run as pi user with bluetooth group access
# For port 80, change to User=root and Group=root
User=pi
Group=pi
SupplementaryGroups=bluetooth

Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=/home/pi/pi-bt-hub/src
Environment=BT_HUB_DB_PATH=/home/pi/pi-bt-hub/data/bt_hub.db
Environment=BT_HUB_LOG_LEVEL=INFO

# Bridge integration disabled by default
# To enable, run: sudo systemctl edit bt-hub
# And add:
#   [Unit]
#   Wants=bt-bridge.service
#   After=bt-bridge.service
#   [Service]
#   Environment=BT_HUB_BRIDGE_ENABLED=true

Environment=BT_HUB_BRIDGE_ENABLED=false
Environment=BT_HUB_BRIDGE_URL=http://localhost:8081

[Install]
WantedBy=multi-user.target
